# files/nitrous_bootstrap.conf
description "Nitrous Bootstrap"

start on (starting ssh)

console output

exec /bin/bash /nitrous/bootstrap.sh
