# files/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/etc/init.d/ssh start

[program:shell]
command=/sbin/init --default-console none
